# playbooks/files/cleanup.service
[Unit]
Description=Server Cleanup Service

[Service]
User=root
ExecStart=/cleanup.sh
Restart=always

[Install]
WantedBy=multi-user.target